# scripts/templates/exam-platform.service.template
[Unit]
Description=Exam Platform Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/exam-platform/backend
Environment=JAVA_HOME=/usr/lib/jvm/java-17-openjdk
Environment=PATH=/usr/lib/jvm/java-17-openjdk/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=SPRING_PROFILES_ACTIVE=mysql
Environment=DB_HOST=127.0.0.1
Environment=DB_PORT=3306
Environment=DB_NAME=exam_platform
Environment=DB_USERNAME=exam_user
Environment=DB_PASSWORD=change_me
Environment=APP_JWT_SECRET=replace_with_secure_secret_min_32_bytes
Environment=APP_CORS_ALLOW_ALL_ORIGINS=false
Environment=APP_SEED_DEMO_DATA_ENABLED=false
Environment=APP_PUBLIC_REGISTRATION_ENABLED=true
Environment=APP_DEMO_PASSWORD_RESET_ENABLED=false
Environment=APP_BOOTSTRAP_ADMIN_USERNAME=
Environment=APP_BOOTSTRAP_ADMIN_PASSWORD=
Environment=APP_AGENT_ENABLED=true
Environment=APP_AGENT_API_KEY=
Environment=APP_AGENT_BASE_URL=https://api.openai.com/v1/chat/completions
Environment=APP_AGENT_MODEL=gpt-4o-mini
Environment=APP_NOTIFY_SMS_WEBHOOK_URL=
Environment=SPRING_MAIL_HOST=
Environment=SPRING_MAIL_PORT=
Environment=SPRING_MAIL_USERNAME=
Environment=SPRING_MAIL_PASSWORD=
# 先在本机执行: cd /opt/exam-platform/backend && mvn -DskipTests package
ExecStart=/usr/lib/jvm/java-17-openjdk/bin/java -jar /opt/exam-platform/backend/target/exam-platform-1.0.0.jar
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
